# 2-trimestre/Bot/script/3/crontab
*/5 * * * * ruby /usr/local/bin/bot14controller
